# {{cookiecutter.project_name}}/varnish.supervisord.conf
[supervisord]
nodaemon=true

[program:varnishd]
command=/usr/sbin/varnishd -f /etc/varnish/default.vcl  -a 0.0.0.0:80 -s malloc,1G -F
autorestart=true